# launchd/com.otsuka.fishing-sns-bot.feature.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.otsuka.fishing-sns-bot.feature</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/koheiotsuka/fishing-sns-bot/scripts/run_feature.sh</string>
    </array>

    <!-- 毎日 19:00 JST に魚種特集投稿を実行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>19</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/Users/koheiotsuka/fishing-sns-bot/data/logs/feature.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/koheiotsuka/fishing-sns-bot/data/logs/feature.err.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/koheiotsuka/fishing-sns-bot</string>
</dict>
</plist>
